# com.nerijus.news-digest-enricher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.nerijus.news-digest-enricher</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/nerijus/Python/masikonis/news-digest/langchain_env/bin/python3</string>
        <string>/Users/nerijus/Python/masikonis/news-digest/src/content_enricher.py</string>
        <string>--config</string>
        <string>/Users/nerijus/Python/masikonis/news-digest/src/config.json</string>
    </array>
    <key>StartInterval</key>
    <integer>3600</integer>
    <key>StandardOutPath</key>
    <string>/Users/nerijus/Python/masikonis/news-digest/output.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/nerijus/Python/masikonis/news-digest/output.log</string>
</dict>
</plist> 