# example-com.example-news-digest-weekly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.example.news-digest-rss-scraper</string>
    <key>ProgramArguments</key>
    <array>
        <string>/path/to/your/python3</string>
        <string>/path/to/your/rss_scraper.py</string>
        <string>--config</string>
        <string>/path/to/your/config.json</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>0</integer> <!-- 0 = Sunday -->
        <key>Hour</key>
        <integer>19</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/nerijus/Python/masikonis/news-digest/output.log</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <dict>
        <key>PathState</key>
        <dict>
            <key>/Users/nerijus/Python/masikonis/news-digest/src/config.json</key>
            <true/>
        </dict>
    </dict>
</dict>
</plist>
